# .config/systemd/user/mpv-fullscreen.timer
[Unit]
Description=Check for full screen session

[Timer]
OnCalendar=*:*:0/5
AccuracySec=5s
Unit=%i.service

[Install]
WantedBy=timers.target
